# .config/sv/pipewire/run
#!/bin/bash
sleep 2
exec pipewire > /tmp/pipe.log 2>&1
