# .config/sv/xdg-desktop-portal-wlr/run
#!/bin/bash
exec /usr/libexec/xdg-desktop-portal-wlr
