# .config/sv/xdg-desktop-portal/run
#!/bin/bash
sleep 3
if pgrep -x "pipewire" > /dev/null; then
	exec /usr/libexec/xdg-desktop-portal -r > /tmp/xdp.log 2>&1
else
	exit 1
fi
